# InitScripts/atvclient
#!/sbin/openrc-run
# Start ATVclient to use the Apple Remote controller
# By Gulivert - http://gulivert.ch

description="Start and stop ATVClient"
  
start() {
	ebegin "Starting ATVClient"
	start-stop-daemon --start -b -m --pidfile /var/run/atvclient.pid --make-pidfile --exec /usr/local/sbin/atvclient ${ATVCLIENT_OPTS}
	eend $?
}
  
stop() {
	ebegin "Stopping ATVClient"
	start-stop-daemon --stop --pidfile /var/run/atvclient.pid
	eend $?
}

